# testapp2/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:fluxbox]
command=fluxbox -display :99
autorestart=true
priority=20
environment=DISPLAY=":99"
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox.err

[program:x11vnc]
command=x11vnc -display :99 -forever -nopw -shared -rfbport 5900 -xkb
autorestart=true
priority=30
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

[program:reflex]
command=poetry run reflex run --env dev --frontend-port 3000 --backend-port 8000 --loglevel info
directory=/app
autorestart=true
priority=100
environment=DISPLAY=":99"
stdout_logfile=/var/log/reflex.log
stderr_logfile=/var/log/reflex.err
